# target/linux/mediatek/mt7622/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	iptime,ax8004m)
		fw_setenv boot_from firmware2
		;;
	linksys,e8450)
		mtd erase senv || true
		;;
	esac
}
